# files/maintenance/services/foot-server.service
[Unit]
Description=Foot server
[Service]
Type=exec
ExecStart=/usr/bin/foot -s
[Install]
WantedBy=multi-user.target
