# files/maintenance/services/hyprland.service
[Unit]
Description=Hyprland
[Service]
Type=exec
ExecStart=/usr/bin/start-hyprland
[Install]
WantedBy=multi-user.target
