# files/maintenance/services/hypridle.service
[Unit]
Description=Hypridle
[Service]
Type=exec
ExecStart=/usr/bin/hypridle
[Install]
WantedBy=multi-user.target
