# files/maintenance/services/wall-auto.service
[Unit]
Description=Wall-auto
[Service]
Type=exec
ExecStart=/home/sdobrau/bin/wall-auto.sh
[Install]
WantedBy=multi-user.target
