# files/maintenance/services/waybar.service
[Unit]
Description=Waybar
[Service]
Type=exec
ExecStart=/usr/bin/waybar
[Install]
WantedBy=multi-user.target
